# services/realtime-notifications/systemd/bugcatcher-realtime-notifications.service
[Unit]
Description=BugCatcher realtime notifications websocket service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/bugcatcher/services/realtime-notifications
Environment=BUGCATCHER_CONFIG_PATH=/var/www/bugcatcher/shared/config.php
ExecStart=/usr/bin/node /var/www/bugcatcher/services/realtime-notifications/server.cjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
