# services/realtime-notifications/systemd/webtest-realtime-notifications.service
[Unit]
Description=WebTest realtime notifications websocket service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mvine
Group=mvine
WorkingDirectory=/var/www/webtest/services/realtime-notifications
Environment=WEBTEST_CONFIG_PATH=/var/www/webtest/shared/config.php
ExecStart=/usr/bin/node /var/www/webtest/services/realtime-notifications/server.cjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
